# deploy/soundit-psa-technician-queue.service
[Unit]
Description=Sound PSA — AI Technician queue worker
After=network.target mysql.service

[Service]
User=www-data
Restart=always
RestartSec=3
WorkingDirectory=/var/www/psa
ExecStart=/usr/bin/php /var/www/psa/artisan queue:work --queue=technician --sleep=3 --tries=2 --max-time=3600
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
